# config/.local/etc/profile.d/amd-rx6700xt-rocm.sh


# Add for support acceleration on ROCm
# rx6700xt specific
export HSA_OVERRIDE_GFX_VERSION=10.3.0
