# config/.local/etc/profile.d/ccache.sh

# Direct mode for speed up
export SCCACHE_DIRECT=true
# Cache size
export SCCACHE_CACHE_SIZE="25G"