# config/.local/etc/profile.d/user-python.sh

prepend_path "$HOME/.local/opt/python/bin"
